# chish.service
[Unit]
Description=My Service
After=network.target

[Service]
ExecStart=/bin/bash /kragg/innit.sh

[Install]
WantedBy=multi-user.target